# systemd/pyfarm-control.service
[Unit]
Description=PyFarm Control Loop for Mushroom Chamber
After=network.target

[Service]
Type=simple
User=pyfarm
Group=pyfarm
WorkingDirectory=/opt/pyfarm
ExecStart=/usr/local/bin/pyfarm grow start /etc/pyfarm/grow.yaml --api-port 8765 --db /var/lib/pyfarm/pyfarm.db
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# Give GPIO access
SupplementaryGroups=gpio

# Hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/pyfarm
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
